# deployment/systemd/faceit-api.service
[Unit]
Description=Faceit AI Bot Docker Stack
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/faceit-ai-bot
EnvironmentFile=/opt/faceit-ai-bot/.env
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
